# dotfiles/lglogin/Library/LaunchAgents/net.garron.mac.lglogin.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>net.garron.mac.lglogin</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/lgarron/Code/git/github.com/lgarron/dotfiles/setup/scripts/trampoline.fish</string>
        <string>/Users/lgarron/.config/lglogin/lglogin.fish</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <dict>
        <key>Minute</key>
        <integer>2</integer>
    </dict>
</dict>
</plist>
