# dotfiles/obsidian-backup-macos/Library/LaunchAgents/net.garron.mac.dotfiles.mirror.obsidian-backup-macos.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>net.garron.mac.dotfiles.mirror.obsidian-backup-macos</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/trampoline</string>
        <string>/Users/lgarron/Code/git/github.com/lgarron/dotfiles/mirrored/obsidian-backup/obsidian-backup.ts</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
